# deploy/recipe.service
[Unit]
Description=Recipe API Spring Boot Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/recipe-api
ExecStart=/usr/bin/java -jar -Dspring.profiles.active=prod /home/ubuntu/recipe-api/recipe-api-0.0.1-SNAPSHOT.jar
EnvironmentFile=/home/ubuntu/recipe-api/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
